# deploy/netlogs-web.service
[Unit]
Description=NetLogs Web Application (Gunicorn)
Documentation=https://github.com/your-org/netlogs
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=net
Group=net
WorkingDirectory=/home/net/net-logs
Environment="PATH=/home/net/net-logs/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DJANGO_SETTINGS_MODULE=netlogs.settings"

# Gunicorn with optimal settings for production
ExecStart=/home/net/net-logs/venv/bin/gunicorn netlogs.wsgi:application \
    --bind 0.0.0.0:8001 \
    --workers 4 \
    --threads 2 \
    --worker-class gthread \
    --worker-tmp-dir /dev/shm \
    --timeout 30 \
    --keep-alive 5 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --access-logfile - \
    --error-logfile - \
    --capture-output \
    --enable-stdio-inheritance

ExecReload=/bin/kill -s HUP $MAINPID
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM

Restart=always
RestartSec=5

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netlogs-web

[Install]
WantedBy=multi-user.target
