# deploy/zentryc-disk-cleanup.service
[Unit]
Description=Zentryc ClickHouse Disk Cleanup
After=network.target docker.service

[Service]
Type=oneshot
User=net
WorkingDirectory=/home/net/zentryc
ExecStart=/home/net/zentryc/venv/bin/python -m fastapi_app.cli.disk_cleanup --threshold 90 --retention 60
StandardOutput=append:/home/net/zentryc/logs/disk_cleanup.log
StandardError=append:/home/net/zentryc/logs/disk_cleanup.log

# Restart on failure after 5 minutes
Restart=on-failure
RestartSec=300

[Install]
WantedBy=multi-user.target
